# scraper/crontab
0 * * * * python3 scraper/scraper.py >> scraper/cron.log
